# bluetooth_speaker_manager/rootfs/etc/services.d/bluetooth_manager/finish
#!/usr/bin/with-contenv bashio
bashio::log.warning "Bluetooth Speaker Manager s'est arrêté"
